# packaging/systemd/acornops-agentv.service
[Unit]
Description=AcornOps AgentV
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=acornops-agent
Group=acornops-agent
EnvironmentFile=/etc/acornops/agentv.env
WorkingDirectory=/opt/acornops/agentv
ExecStart=/usr/bin/node /opt/acornops/agentv/dist/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/acornops-agentv

[Install]
WantedBy=multi-user.target
